# systemd/app.service
[Unit]
Description=PacketBench
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/packetbench
EnvironmentFile=/opt/packetbench/.env
ExecStart=/opt/packetbench/.venv/bin/python scripts/run.py
User=www-data
Group=www-data
Restart=always
RestartSec=3
RestartPreventExitStatus=3

[Install]
WantedBy=multi-user.target
